# templates/awx-nginx-cert.path.j2
[Path]
PathChanged={{ awx_ssl_crt_path }}

[Install]
WantedBy=multi-user.target
