# templates/awx-nginx-cert.path.j2
[Path]
PathChanged={{ hostvars['localhost']['ssl_crt_path'] }}

[Install]
WantedBy=multi-user.target
