# templates/awx-nginx-cert.path.j2
[Path]
PathChanged={{ hostvars['localhost']['awx_nginx_crt_path'] }}

[Install]
WantedBy=multi-user.target
